# netrt.service.example
# This is an example systemd unit file for running CORRECT directly on the host.
# Update the paths, user, group, and config locations for your deployment.

[Unit]
Description=CORRECT DICOM Processing Application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=appuser
Group=appgroup

# Set the working directory to the repository root.
WorkingDirectory=/opt/CORRECT

# Ensure the configured directories exist and are writable by User/Group:
#   directories.working
#   directories.logs
#   reporting.output_directory
#   anonymization.pid_mapping_file parent directory
Environment="PYTHONUNBUFFERED=1"

# Command to start the application.
# Point --config at the same YAML file that defines the listener, destination,
# working/log/report directories, and feature flags.
ExecStart=/usr/bin/python3 /opt/CORRECT/main.py --config /etc/correct/config.yaml

# Restart policy
Restart=on-failure
RestartSec=5s

# Standard output and error logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
